# data/uploadfile.conf
start on startup
task
exec /etc/init.d/uploadfile
